# deploy/pkg/debian/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable ai-agent-vision.service
systemctl restart ai-agent-vision.service
echo "Post-installation script executed successfully."
