# deploy/pkg/debian/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable ai-agent-vision.service
systemctl restart ai-agent-vision.service

# Only copy .env if it does not already exist
ENV_PATH="/opt/ai-agent-vision/conf/.env"
ENV_EXAMPLE_PATH="/opt/ai-agent-vision/conf/.env.example"
if [ ! -f "$ENV_PATH" ] && [ -f "$ENV_EXAMPLE_PATH" ]; then
    cp "$ENV_EXAMPLE_PATH" "$ENV_PATH"
    echo ".env file created from .env.example. Please review and update as needed."
elif [ -f "$ENV_PATH" ]; then
    echo ".env file already exists and has not been updated."
fi

echo "Post-installation script executed successfully."
