# deploy/pkg/debian/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable ai-agent-vision.service
